# deploy/systemd/compras-gov-crawl.service
# /etc/systemd/system/compras-gov-crawl.service
# ComprasGov crawler — daily at 07:00 UTC

[Unit]
Description=ComprasGov Crawler — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/usr/bin/python3 /opt/extra-consultoria/scripts/crawl/monitor.py --source compras_gov --mode full
User=extra-consultoria
StandardOutput=journal
StandardError=journal
