# deploy/systemd/coverage-report.service
# /etc/systemd/system/coverage-report.service
# Daily coverage report + snapshot at 09:00 UTC

[Unit]
Description=Coverage Gap Report + Snapshot — Extra Consultoria
After=pncp-crawl-full.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStartPre=/usr/bin/python3 /opt/extra-consultoria/scripts/crawl/monitor.py --report-coverage
ExecStart=/opt/extra-consultoria/.venv/bin/python -m scripts.local_datalake coverage --snapshot --export
User=extra-consultoria
StandardOutput=journal
StandardError=journal
