# deploy/systemd/coverage-report-weekly.timer
# /etc/systemd/system/coverage-report-weekly.timer
# Weekly coverage report — every Monday at 08:00 UTC

[Unit]
Description=Weekly Coverage Report Timer — Extra Consultoria
Requires=coverage-report-weekly.service

[Timer]
OnCalendar=Mon *-*-* 08:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
