# deploy/systemd/coverage-report-weekly.service
# /etc/systemd/system/coverage-report-weekly.service
# Weekly coverage report — every Monday at 08:00 UTC

[Unit]
Description=Weekly Coverage Report — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/opt/extra-consultoria/.venv/bin/python -m scripts.reports.coverage_weekly
User=extra-consultoria
StandardOutput=journal
StandardError=journal
