# deploy/systemd/extra-check-alerts.service
# /etc/systemd/system/extra-check-alerts.service
# Periodic alert checking for system health
# Extra Consultoria — Story TD-5.5

[Unit]
Description=Extra Consultoria — Check System Alerts
Documentation=https://github.com/extra-consultoria/docs/ops/monitoring.md
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/usr/bin/python3 /opt/extra-consultoria/scripts/check-alerts.py --json
User=extra-consultoria
StandardOutput=journal
StandardError=journal
