# deploy/systemd/extra-crawl-selenium.service
# /etc/systemd/system/extra-crawl-selenium.service
# Selenium crawler — INTENTIONALLY DISABLED
#
# HISTORICO:
#   Selenium foi removido como fonte de dados independente (Story 1.5 fix).
#   Commit: 525ff22 (retroactive audit remediation)
#   Referencia: registry.py linha 352-354:
#     "selenium REMOVED as a source (Story 1.5 fix).
#      Selenium e um metodo de crawl, nao uma fonte de dados.
#      Fontes que usam selenium declararam o modo 'selenium' em seus modos suportados."
#
# COBERTURA ATUAL:
#   Portais JS-rendered sao crawleados via adaptadores que usam selenium
#   internamente como metodo (ex: transparencia_crawler com modo "selenium").
#   Nao ha necessidade de um servico separado para selenium.
#
# REATIVACAO:
#   Se um novo mecanismo de crawl selenium independente for necessario,
#   registrar como source em registry.py e atualizar este servico.
#   Ate la, este servico falha explicitamente para evitar execucao silenciosa.

[Unit]
Description=Extra Consultoria — Selenium JS Portal Crawler [DISABLED]
After=network-online.target
Wants=network-online.target
OnFailure=onfailure@%n.service

[Service]
Type=oneshot
ExecStartPre=/bin/echo "DISABLED: selenium nao e uma fonte registrada (Story 1.5 fix)"
ExecStart=/bin/false
User=extra-consultoria
StandardOutput=journal
StandardError=journal
