# deploy/systemd/onfailure@.service
# /etc/systemd/system/onfailure@.service
# Template unit for failure notifications via webhook
# Required: add WEBHOOK_URL to /opt/extra-consultoria/.env
#
# Usage: referenced by OnFailure=onfailure@%n.service in other units
# Example WEBHOOK_URL: https://hooks.slack.com/services/...

[Unit]
Description=OnFailure alert for %i

[Service]
Type=oneshot
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/usr/bin/curl -sf --max-time 10 -X POST -H "Content-Type: application/json" -d '{"service":"%i","host":"%H","status":"failed"}' "${WEBHOOK_URL}"
