# deploy/systemd/pncp-crawl-full.timer
# /etc/systemd/system/pncp-crawl-full.timer
# Full PNCP crawl — daily at 05:00 UTC (02:00 BRT)

[Unit]
Description=Daily PNCP Full Crawl Timer — Extra Consultoria
Requires=pncp-crawl-full.service

[Timer]
OnCalendar=*-*-* 05:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
