# deploy/systemd/pncp-contracts.service
# /etc/systemd/system/pncp-contracts.service
# PNCP contracts crawler / incremental — project venv required

[Unit]
Description=PNCP Contracts Crawler — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/usr/bin/flock --nonblock /run/lock/pncp-contracts.lock \
  /opt/extra-consultoria/.venv/bin/python -m scripts.crawl.run_contracts_incremental --days 7
User=extra-consultoria
Group=extra-consultoria
TimeoutStartSec=90min
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pncp-contracts

[Install]
WantedBy=multi-user.target
