# deploy/systemd/pncp-crawl-inc.service
# /etc/systemd/system/pncp-crawl-inc.service
# Incremental PNCP crawl — uses project venv (psycopg2 etc.)

[Unit]
Description=PNCP DataLake Incremental Crawl — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
# flock prevents concurrent writers
ExecStart=/usr/bin/flock --nonblock /run/lock/pncp-crawl-inc.lock \
  /opt/extra-consultoria/.venv/bin/python -m scripts.crawl.monitor --source pncp --mode incremental
User=extra-consultoria
Group=extra-consultoria
TimeoutStartSec=45min
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pncp-crawl-inc

[Install]
WantedBy=multi-user.target
