# deploy/systemd/pncp-crawl-full.service
# /etc/systemd/system/pncp-crawl-full.service
# Full crawl: all UFs × all modalidades, daily at 05:00 UTC

[Unit]
Description=PNCP DataLake Full Crawl — Extra Consultoria
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/opt/extra-consultoria/.venv/bin/python /opt/extra-consultoria/scripts/crawl/monitor.py --source pncp --mode full
User=extra-consultoria
StandardOutput=journal
StandardError=journal
