# deploy/systemd/sc-compras-crawl.service
# /etc/systemd/system/sc-compras-crawl.service
# SC Compras crawler — incremental semanal (domingo 09:00 UTC)

[Unit]
Description=SC Compras Crawler — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/opt/extra-consultoria/.venv/bin/python /opt/extra-consultoria/scripts/crawl/monitor.py --source sc_compras --mode incremental
User=extra-consultoria
StandardOutput=journal
StandardError=journal
