# deploy/systemd/transparencia-crawl.service
# /etc/systemd/system/transparencia-crawl.service
# Transparencia crawler — weekly (Sun 06:00 UTC)
# Template-driven scraping de portais de transparencia municipais

[Unit]
Description=Transparencia Crawler — Extra Consultoria
After=network-online.target
Wants=network-online.target
OnFailure=extra-onfailure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/extra-consultoria
Environment=PYTHONPATH=/opt/extra-consultoria
EnvironmentFile=/opt/extra-consultoria/.env
ExecStart=/usr/bin/python3 /opt/extra-consultoria/scripts/crawl/monitor.py --source transparencia --mode full
User=extra-consultoria
StandardOutput=journal
StandardError=journal
